# osquery/roles/fleet/templates/kolide-fleet-proxy.socket.j2
[Socket]
ListenStream=0.0.0.0:443
Transparent=yes

[Install]
WantedBy=sockets.target
